# systemd/hosaka-webserver.service
[Unit]
Description=Hosaka Web Server (FastAPI — serves SPA + API)
# Default bind is loopback for the appliance's local kiosk path, but operators
# can override HOSAKA_WEB_HOST=0.0.0.0 in /etc/hosaka/env when they want the
# node to be reachable over Tailscale/LAN for Hosaka peer discovery.
# We still don't need network-online.target here; the service can come up
# before the network does and start answering once the interface arrives.
After=local-fs.target
After=hosaka-mode.service

[Service]
Type=simple
WorkingDirectory=/opt/hosaka-field-terminal
Environment=HOSAKA_WEB_HOST=127.0.0.1
EnvironmentFile=-/etc/hosaka/env
Environment=HOSAKA_STATE_PATH=/var/lib/hosaka/state.json
Environment=HOSAKA_BOOT_MODE=kiosk
Environment=PYTHONUNBUFFERED=1
# Surface missing required env vars (e.g. OPENAI_API_KEY) in the journal
# at boot rather than as a vague 502 from the first /api/v1/voice request.
ExecStartPre=/opt/hosaka-field-terminal/.venv/bin/python -m hosaka.secrets check --required-only --no-probe
ExecStart=/opt/hosaka-field-terminal/.venv/bin/uvicorn \
    hosaka.web.server:app \
    --host ${HOSAKA_WEB_HOST} \
    --port 8421 \
    --workers 1 \
    --log-level info
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Resource caps so a runaway request handler can't OOM the whole Pi or
# starve the kiosk of CPU. Numbers are sized for a Pi 3B (1 GB RAM, 4 cores).
MemoryMax=300M
MemoryHigh=240M
CPUQuota=75%
TasksMax=64

[Install]
WantedBy=multi-user.target
