# systemd/hosaka-oom-guard.service
[Unit]
Description=Hosaka memory pressure watchdog (kills kiosk before swap thrash takes the SD card down)
DefaultDependencies=no
After=basic.target
Before=shutdown.target
Conflicts=shutdown.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hosaka-oom-guard
Restart=always
RestartSec=5

# Resource caps: this is the LAST process we want hogging RAM. Cap hard.
MemoryHigh=20M
MemoryMax=40M
CPUQuota=15%
TasksMax=8

# Make sure the OOM killer takes anything else first (-1000 = kernel never picks).
OOMScoreAdjust=-900
Nice=-5

# Run as root so we can write /proc/sys/vm/drop_caches and pkill across users.
User=root

# The script writes /var/log/hosaka/oom-guard.log directly. Send stdout to
# the journal (and out, since systemd will pick it up) — do NOT also append
# to the same file or every line ends up duplicated.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
