# systemd/picoclaw-gateway.service
[Unit]
Description=Picoclaw Agent Gateway
# No network-online dependency: the gateway reaches out to the LLM provider
# lazily on first request, so a slow NIC bring-up shouldn't block boot.
After=local-fs.target
After=hosaka-mode.service
# Don't bring up picoclaw at all if the Pi booted into build mode.
ConditionPathExists=!/boot/firmware/hosaka-build-mode

[Service]
Type=simple
# Run as the user who owns the picoclaw config (~/.picoclaw).
# Override with: sudo systemctl edit picoclaw-gateway --force
User=operator
ExecStart=/usr/local/bin/picoclaw gateway
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# This process has been observed to sustain ~14% of one core at idle
# (likely a tight poll loop). Cap it so it can't dominate the Pi 3B's
# limited CPU budget while we figure out a real fix upstream.
CPUQuota=25%
Nice=10
MemoryMax=200M

[Install]
WantedBy=multi-user.target
