# systemd/dropins/ssh.service.d/oom-protection.conf
# SSH guard dog.
#
# When the Pi 3B starts thrashing swap (npm build + chromium + cursor server
# all at once), Linux's OOM killer picks "the worst offender" — historically
# that's been sshd because it has lots of children + a high RSS sum. Losing
# sshd means we lose remote access right when we need it most to recover.
#
# OOMScoreAdjust=-1000 is the strongest possible bias: the OOM killer will
# only touch sshd if literally nothing else is killable. Combined with the
# bias-up we apply to chromium / node in the kiosk autostart and the
# `hosaka` CLI, the kill order under pressure becomes:
#
#   chromium  >  node/npm  >  picoclaw-gateway  >  uvicorn  >>  sshd
#
# This file is a systemd drop-in; install it at:
#   /etc/systemd/system/ssh.service.d/oom-protection.conf
# or use scripts/install_hosaka_lean.sh which handles it.

[Service]
OOMScoreAdjust=-1000
# Don't let sshd itself ever be cgroup-OOM-killed via the kernel "kill the
# whole unit" knob — even if a single forked session goes wild, we want
# the parent sshd to keep listening.
OOMPolicy=continue
